# templates/supervisor/encoder.ini
[program:encoder]
command=php artisan queue:work --queue=encoding_sv_{server_id} --timeout=12500 --tries=1
directory=/var/www/{server_name}/
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/www/{server_name}/encoder.out.log
stderr_logfile=/var/www/{server_name}/encoder.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=www
priority=999
numprocs=5
process_name=%(program_name)s_%(process_num)02d